# yocto_template/files/server.service
[Unit]
Description=GRPC ##SERVICE_NAME Server

[Service]
Type=simple
ExecStart=/usr/bin/##SERVICE_NAME

[Install]
WantedBy=multi-user.target
